# deploy/systemd/ngo-ai-omlx-tunnel.service
[Unit]
Description=SSH tunnel to Mac Studio oMLX for NGO AI Demo
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=radek
ExecStart=/usr/bin/ssh -N -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -L 127.0.0.1:18585:127.0.0.1:8585 mac-studio-ngo-demo
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
